# src/sysmaint/data/systemd/sysmaint-pihole.timer
[Unit]
Description=Run sysmaint Pi-hole update monthly
Documentation=https://github.com/Stephen-Kennedy/auto_updater_package

[Timer]
# First Sunday of each month at 04:00 + 1h jitter.
OnCalendar=Sun *-*-1..7 04:00:00
RandomizedDelaySec=1h
Persistent=true
Unit=sysmaint-pihole.service

[Install]
WantedBy=timers.target
